# parking-camera.service
[Unit]
Description=Parking Camera Client Service
Documentation=https://github.com/your-org/parking-client
After=network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/parking-client
ExecStart=/opt/parking-client/venv/bin/python /opt/parking-client/pi_camera_client.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="CAMERA_CONFIG=/opt/parking-client/config.json"

# security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/parking-client
ReadWritePaths=/tmp
ReadWritePaths=/var/log/parking-client

# resource limits
MemoryMax=512M
CPUQuota=80%

# logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=parking-camera

[Install]
WantedBy=multi-user.target